# data/com.leslie.pulse.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.leslie.pulse.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/leslie/Workspace/ainvest/Leslie-stock/data/daily_pulse.sh</string>
  </array>

  <!-- 每天 18:00（美股盘中、A 股已收）触发；周末也跑（拿不到新数据也无害） -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/leslie/Workspace/ainvest/Leslie-stock/data/logs/launchd-out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/leslie/Workspace/ainvest/Leslie-stock/data/logs/launchd-err.log</string>

  <!-- 错过时段（电脑关机等）启动后自动补跑 -->
  <key>RunAtLoad</key>
  <false/>
  <key>StartOnMount</key>
  <false/>
</dict>
</plist>
