# scripts/com.leslie.stock.premarket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leslie.stock.premarket</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/leslie/Workspace/ainvest/Leslie-stock/scripts/premarket_job.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leslie/Workspace/ainvest/Leslie-stock/logs/premarket_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leslie/Workspace/ainvest/Leslie-stock/logs/premarket_launchd.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/leslie/Workspace/ainvest/Leslie-stock</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/leslie/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
